# systemd/fantasy-espn.timer
[Unit]
Description=Run Fantasy Baseball ESPN forecaster sync every Sunday at 5:30 AM

[Timer]
OnCalendar=Sun *-*-* 05:30:00
AccuracySec=5m
Persistent=true

[Install]
WantedBy=timers.target
