# systemd/fantasy-run-all.service
[Unit]
Description=Fantasy Baseball full sync and report pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/rlm/src/fantasy
Environment=PATH=/home/rlm/.local/bin:/usr/local/bin:/usr/bin:/usr/sbin
ExecStart=fantasy-run-all
StandardOutput=journal
StandardError=journal
