# systemd/fantasy-sync.timer
[Unit]
Description=Run Fantasy Baseball Yahoo sync daily at 6:00 AM

[Timer]
OnCalendar=*-*-* 06:00:00
AccuracySec=5m
Persistent=true

[Install]
WantedBy=timers.target
