# systemd/fantasy-report.timer
[Unit]
Description=Email Fantasy Baseball pitcher start reports every Sunday at 6:00 AM

[Timer]
OnCalendar=Sun *-*-* 06:00:00
AccuracySec=5m
Persistent=true

[Install]
WantedBy=timers.target
